# services/app_v1/app_v1.service
[Unit]
Description=RestAPI app v1
After=network.target

[Service]
ExecStart=/usr/bin/app_v1
Restart=always

[Install]
WantedBy=multi-user.target
